# systemd/clawdbot-watchdog.timer
[Unit]
Description=ClawdBot Gateway Watchdog Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
