# examples/watchdog.service.example
# ClawdBot Watchdog — oneshot health check service
# Copy to ~/.config/systemd/user/clawdbot-watchdog.service
# Triggered by clawdbot-watchdog.timer

[Unit]
Description=ClawdBot Gateway Watchdog Health Check

[Service]
Type=oneshot
ExecStart=/root/.clawdbot/scripts/gateway-watchdog.sh
